# systemd/system/auto-update-discord.timer
[Unit]
Description=Run the auto-update-discord.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
