# systemd/system/auto-update-discord.service
[Unit]
Description=Automatically Update Discord

[Service]
Type=simple
ExecStart=/usr/local/bin/update_discord.sh
